# deploy/systemd/palworld.service
[Unit]
Description=Palworld Dedicated Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=palworld
Group=palworld
WorkingDirectory=/opt/palworld
ExecStart=/opt/palworld/PalServer.sh -useperfthreads -NoAsyncLoadingThread -UseMultithreadForDS
# Let Palworld flush its save on stop instead of being SIGKILLed.
KillSignal=SIGINT
TimeoutStopSec=90
Restart=on-failure
RestartSec=10
# Memory-leak guard (Palworld climbs 15-25 GB over time). Sized for a
# 16 GiB CT; raise on a bigger box. OOMPolicy=stop lets Restart bring it back.
MemoryHigh=22G
MemoryMax=24G
OOMPolicy=stop
LimitNOFILE=100000
SyslogIdentifier=palworld

[Install]
WantedBy=multi-user.target
