# deploy/systemd/rallypoint-cmd.service
[Unit]
Description=Palworld Management Panel
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=palworld
Group=palworld
WorkingDirectory=/opt/rallypoint-cmd
EnvironmentFile=/etc/rallypoint-cmd/panel.env
ExecStart=/usr/bin/node /opt/rallypoint-cmd/apps/server/dist/server.js
Restart=always
RestartSec=3
MemoryMax=512M
SyslogIdentifier=rallypoint-cmd
# IMPORTANT: do NOT set NoNewPrivileges=true here — the panel calls
# `sudo systemctl` to control palworld.service, and NoNewPrivileges
# breaks sudo. (If you want that hardening, switch service control to a
# polkit rule first; see deploy/README.md.)

[Install]
WantedBy=multi-user.target
